# letsencrypt/root.crontab
* * * * * echo "Refresh LetsEncrypt certificate here." && date > /usr/local/apache2/htdocs/.well-known/acme-challenge
